# .config/systemd/user/network-check.timer
[Unit]
Description=Run network check every minute

[Timer]
OnBootSec=30sec
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target

